# core/supervisor/prod.conf
[program:webserver]
directory = /opt/mis
command = gunicorn project.wsgi --bind=0.0.0.0:80 --workers=5

[program:check_tasks]
directory = /opt/mis
command = python3 manage.py task_checker

